# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user = worker
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:failed_email]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "while [ true ]; do (php /var/www/failed_email.php &); sleep 60; done"
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true

[program:worker-1]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/consumer.php
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker-2]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/consumer.php
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true